# mumble/etc/supervisor/conf.d/murmurd.conf
[program:murmurd]
command=/usr/sbin/murmurd -fg -ini /data/mumble-server.ini
